# .github/celery-worker.service
[Unit]
Description=Celery Worker for DRF Vue Backend
After=network.target redis-server.service mysql.service
Requires=redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/DRF_VUE/drf_vue/backend
Environment="DJANGO_ENV=production"
Environment="PATH=/home/DRF_VUE/drf_vue/backend/venv/bin"
ExecStart=/home/DRF_VUE/drf_vue/backend/venv/bin/celery -A config.celery worker \
    --concurrency=2 \
    --prefetch-multiplier=2 \
    --max-tasks-per-child=100 \
    --loglevel=INFO \
    --logfile=/home/DRF_VUE/drf_vue/backend/logs/celery.log
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=always
RestartSec=10s

# Security settings
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
