# .github/celery-beat.service
[Unit]
Description=Celery Beat for DRF Vue Backend
After=network.target redis-server.service mysql.service
Requires=redis-server.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/home/DRF_VUE/drf_vue/backend
Environment="DJANGO_ENV=production"
Environment="PATH=/home/DRF_VUE/drf_vue/backend/venv/bin"
ExecStart=/home/DRF_VUE/drf_vue/backend/venv/bin/celery -A config.celery beat -l info --schedule=/tmp/celerybeat-schedule
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
Restart=always
RestartSec=10s
StandardOutput=append:/home/DRF_VUE/drf_vue/backend/logs/celery-beat.log
StandardError=append:/home/DRF_VUE/drf_vue/backend/logs/celery-beat-error.log

# 安全设置
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
