# backend/supervisor/celery.conf
; Celery Worker 配置
[program:celery-worker]
command=/path/to/venv/bin/celery -A config.celery worker -l info
directory=/path/to/backend
user=www-data
group=www-data
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
stopsignal=TERM
stdout_logfile=/var/log/celery/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/celery/worker_error.log
environment=DJANGO_ENV="production"

; Celery Beat 配置（定时任务调度器）
[program:celery-beat]
command=/path/to/venv/bin/celery -A config.celery beat -l info --schedule=/var/lib/celery/celerybeat-schedule
directory=/path/to/backend
user=www-data
group=www-data
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=300
stopsignal=TERM
stdout_logfile=/var/log/celery/beat.log
stderr_logfile=/var/log/celery/beat_error.log
environment=DJANGO_ENV="production"
